# indimail-mta-2.8/debian/indimail-mta.postinst.in
#!/bin/sh
# postinst script for indimail
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

prefix=@prefix@
qmaildir=@qmaildir@
mandir=@mandir@
shareddir=@shareddir@
libexecdir=@libexecdir@
qcount=5
qbase=${qmaildir}/queue
mbase=@mbase@
logdir=@logdir@
servicedir=@servicedir@
qsysconfdir=@qsysconfdir@
sysconfdir=@sysconfdir@
cp=/bin/cp
rm=/bin/rm
mv=/bin/mv
####################################
nodksignature=0
notlsarr=0
tcpserver_plugin=1
dkimkeyfn=default

setperm()
{
	local user="$1"; shift
	local group="$1"; shift
	local mode="$1"; shift
	local file="$1"; shift
	/usr/sbin/dpkg-statoverride --list "$file" >/dev/null && return 0
	/usr/sbin/dpkg-statoverride --update --add "$user" "$group" "$mode" "$file"
}

fixperms()
{
	if [ $# -eq 3 ] ; then # directory
	#%dir %attr(0555,root,indimail) /var/indimail/sbin
		lin=$2
		file=$3
		if [ ! -d $file ] ; then
			return 0
		fi
	elif [ $# = 2 ] ; then # file
	#%attr(0555,root,qmail) /var/indimail/bin/rsmtp
		lin=$1
		file=$2
		if [ ! -f $file ] ; then
			return 0
		fi
	fi
	perm=`echo $lin|cut -d\( -f2 | cut -d, -f1`
	perm_f1=`echo $perm | cut -c1`
	own=`echo $lin|cut -d\( -f2 | cut -d, -f2`
	grp=`echo $lin|cut -d\( -f2 | cut -d, -f3|cut -d\) -f1`
	if [ -f $file ] ; then
		if [ $perm_f1 -eq 0 ] ; then # 4755, 0644
			/bin/chown $own:$grp $file
			/bin/chmod $perm $file
		else
			echo $file
			if [ -f /usr/sbin/dpkg-statoverride ] ; then
				setperm $own $grp $perm $file
			else
				/bin/chown $own:$grp $file
				/bin/chmod $perm $file
			fi
		fi
	elif [ -d $file ] ; then
		/bin/chown $own:$grp $file
		/bin/chmod $perm $file
	fi
}

ID=`id -u`
if [ $ID -ne 0 ] ; then
	echo "You are not root" >&2
	exit 1
fi

case "$1" in
    configure)
		if [ " $2" = " " ] ; then
		echo "installing indimail-mta"
		else
		echo "upgrading indimail-mta $2"
		fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
	exit 0
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

#
# We need to do this for debian. SIGH!!
#
if [ -f ${qsysconfdir}/perm_list.q ] ; then
	echo "Fixing Permissions (debian stupidity)"
(
/bin/grep -v "^#" ${qsysconfdir}/perm_list.q | while read line
do
	if [ " $line" = " " ] ; then
		continue
	fi
	fixperms $line
done
) > ${qsysconfdir}/stat.override
fi

# indimail-mta.upgrade/indimail-mta.install are created
# by the preinst script
if [ -f /tmp/indimail-mta.upgrade ] ; then # upgrade
  if [ -f ${libexecdir}/qupgrade.sh ] ; then
    echo "Running Custom Upgrade Script for post upgrade"
    /bin/sh ${libexecdir}/qupgrade.sh post upgrade @version@ $* || true
  fi
  if [ -f /etc/init.d/svscan ] ; then
    cmp -s ${prefix}/bin/qmailctl /etc/init.d/svscan >/dev/null 2>&1 || \
		${cp} ${prefix}/bin/qmailctl /etc/init.d/svscan || true
  fi
  if [ -x /bin/systemctl -a -d /lib/systemd/system ] ; then
    cmp -s ${shareddir}/boot/systemd /lib/systemd/system/svscan.service >/dev/null 2>&1 || true
    ${cp} ${shareddir}/boot/systemd /lib/systemd/system/svscan.service
    /bin/systemctl daemon-reload
  elif [ -x /bin/systemctl -a -d /usr/lib/systemd/system ] ; then
    cmp -s ${shareddir}/boot/systemd /usr/lib/systemd/system/svscan.service >/dev/null 2>&1 || true
    ${cp} ${shareddir}/boot/systemd /usr/lib/systemd/system/svscan.service
    /bin/systemctl daemon-reload
  elif [ -f /sbin/initctl -a -d /etc/init ] ; then
    cmp -s ${shareddir}/boot/upstart /etc/init/svscan.conf >/dev/null 2>&1 || true
    ${cp} ${shareddir}/boot/upstart /etc/init/svscan.conf
  elif [ -d /etc/event.d ] ; then
    cmp -s ${shareddir}/boot/upstart /etc/event.d/upstart >/dev/null 2>&1 || true
    ${cp} ${shareddir}/boot/upstart /etc/event.d
  elif [ -d /System/Library/LaunchDaemons ] ; then
    cmp -s ${shareddir}/boot/indimail.plist /System/Library/LaunchDaemons >/dev/null 2>&1 || true
    ${cp} ${shareddir}/boot/indimail.plist /System/Library/LaunchDaemons
  fi
  (
  # refresh indimail-mta services
  svc_list=""
  for i in clamd greylist.1999 qmail-logfifo qmail-qmqpd.628 \
    qmail-qmtpd.209 qmail-smtpd.25 qmail-smtpd.366 \
    qmail-smtpd.465 qmail-smtpd.587 qmail-daned.1998 \
    qscanq udplogger.3000 libwatch
  do
    if [ -d ${servicedir}/$i ] ; then
      if [ -z "$svc_list" ] ; then
        svc_list="${servicedir}/$i"
      else
        svc_list="$svc_list ${servicedir}/$i"
      fi
    fi
  done
  svc_list="$svc_list ${qsysconfdir}/control/defaultqueue"
  ${prefix}/sbin/svctool --servicedir=${servicedir} --refreshsvc="$svc_list"
  ) >> /tmp/indimail-mta-setup.log
  ${rm} -f /tmp/indimail-mta.upgrade
  exit 0
fi
${rm} -f /tmp/indimail-mta.install
if [ ! "${prefix}" = "/usr" ] ; then
	# Recreate ld.so links and cache
	echo "recreating ld.so cache"
	if [ -d ${sysconfdir}/ld.so.conf.d ] ; then
	(
		if [ -d /usr/lib64 ] ; then
			echo /usr/lib64
		else
			echo /usr/lib
		fi
	) > ${sysconfdir}/ld.so.conf.d/indimail-mta-`uname -m`.conf
	/sbin/ldconfig
	fi
fi

echo "Doing post installation activities for the following"
echo ""
echo " 1. Configure ${logdir} for multilog"
echo " 2. Configure ${servicedir}"
echo " 3. Configure svscanlog service"
echo " 4. Configure standard catch-all accounts, default qmail configuration"
echo " 5. Configure DKIM, Domainkeys signature"
echo " 6. Configure QHPSI for inline virus scanning"
echo " 7. Configure SMTP services on port in 465 25 587"
echo " 8. Configure default queue configuration for sendmail, qmail-inject"
echo " 9. Configure ODMR service"
echo "10. Configure greylisting service"
echo "11. Configure QMTP service"
echo "12. Configure QMQP service"
echo "13. Configure udplogger service"
echo "14. Configure qscanq/clamd/freshclam service"
echo "15. Configure ${qsysconfdir}/control/signatures"
echo "16. Configure tcprules database for SMTP, QMTP, QMQP"
echo "17. Configure indimail-mta startup"
echo ""

echo "Post Installation..."
(
echo "Creating ${logdir}"
if [ ! -d ${logdir} ] ; then
	/bin/mkdir ${logdir}
fi
/bin/chown -R qmaill:nofiles ${logdir}

# fifolog service
${prefix}/sbin/svctool --fifologger=/tmp/logfifo --servicedir=${servicedir}

# Misc/Configuration
${prefix}/sbin/svctool --svscanlog --servicedir=${servicedir} --scanint=60

echo "Creating default qmail configuraton and catch-all alias (${qmaildir}/alias/Maildir) for all system users"
${prefix}/sbin/svctool --config=qmail --postmaster=${qmaildir}/alias/Maildir/

if [ $nodksignature -eq 0 ] ; then
	if [ -x ${prefix}/bin/dknewkey ] ; then
		ver_opt="dkim"
		sign_opt="dkim"
		if [ " $key_bit" = " " ] ; then
			key_bit=1024
		fi
		${prefix}/bin/dknewkey ${qsysconfdir}/control/domainkeys/$dkimkeyfn $key_bit
		/bin/chown indimail:qmail ${qsysconfdir}/control/domainkeys/$dkimkeyfn
		/bin/chmod 640 ${qsysconfdir}/control/domainkeys/$dkimkeyfn
	else
		ver_opt="none"
		sign_opt="none"
	fi
else
	ver_opt="none"
	sign_opt="none"
fi
a_arch=`uname -m`
if [ " $a_arch" = " x86_64" ] ; then
	smtp_soft_mem=536870912
	qmtp_soft_mem=104857600
	qmqp_soft_mem=104857600
else
	smtp_soft_mem=52428800
	qmtp_soft_mem=52428800
	qmqp_soft_mem=52428800
fi

echo "Checking if clamav is installed"
if [ -f /usr/sbin/clamd -a -f /usr/bin/clamdscan ] ; then
	have_clamav=1
	clamdPrefix="/usr"
	if [ -d ${sysconfdir}/clamav ] ; then
		mysysconfdir=${sysconfdir}/clamav
	elif [ -d ${sysconfdir}/clamd.d ] ; then
		mysysconfdir=${sysconfdir}/clamd.d
	elif [ -d ${qsysconfdir}] ; then
		mysysconfdir=${qsysconfdir}
	else
		mysysconfdir=${sysconfdir}
	fi
	qhpsi="/usr/bin/clamdscan %s --config=${mysysconfdir}/scan.conf --fdpass --quiet --no-summary"
else
	have_clamav=0
fi

# SMTP
for port in 465 25 587
do
	extra_opt=""
	if [ $port -eq 465 ] ; then
		extra_opt="--skipsend --ssl"
		extra_opt="$extra_opt --rbl=-rzen.spamhaus.org --rbl=-rdnsbl-1.uceprotect.net"
	elif [ $port -eq 587 ] ; then
		# local authenticated users
		# why bother with spamfilter
		# if you do, then you have a serious
		# problem with your organization, not 
		# with IndiMail
		extra_opt="--skipsend --authsmtp --antispoof --forcetls"
	else
		extra_opt="--remote-authsmtp=plain --localfilter --remotefilter"
		extra_opt="$extra_opt --deliverylimit-count=-1 --deliverylimit-size=-1"
		extra_opt="$extra_opt --rbl=-rzen.spamhaus.org --rbl=-rdnsbl-1.uceprotect.net"
	fi
	if [ $tcpserver_plugin -eq 1 ] ; then
		extra_opt="$extra_opt --shared-objects=1 --use-dlmopen=1"
	fi
	if [ $have_clamav -eq 1 ] ; then
		extra_opt="$extra_opt --qhpsi=\"$qhpsi\""
	fi
	eval ${prefix}/sbin/svctool --smtp=$port --servicedir=${servicedir} \
		--qbase=${qbase} --qcount=${qcount} --qstart=1 \
		--query-cache --dnscheck --password-cache \
		--cntrldir=control --localip=0 --maxdaemons=75 --maxperip=25 --persistdb \
		--starttls --fsync --syncdir --memory=$smtp_soft_mem \
		--chkrecipient --chkrelay --masquerade \
		--min-free=52428800 --content-filter --virus-filter \
		--dmasquerade --infifo=infifo \
		--dkverify=${ver_opt} \
		--dksign=${sign_opt} --private_key=${qsysconfdir}/control/domainkeys/%/$dkimkeyfn \
		$extra_opt
done

# queue parameters in control/defaultqueue for qmail-inject, sendmail
extra_opt=""
if [ $have_clamav -eq 1 ] ; then
	extra_opt="--qhpsi=\"$qhpsi\""
fi
eval ${prefix}/sbin/svctool --queueParam=defaultqueue \
	--qbase=${qbase} --qcount=${qcount} --qstart=1 \
	--min-free=52428800 --fsync --syncdir --virus-filter \
	--dkverify="none" --dksign=$sign_opt \
	--private_key=${qsysconfdir}/control/domainkeys/%/$dkimkeyfn \
	$extra_opt

#odmr
${prefix}/sbin/svctool --smtp=366 --odmr --servicedir=${servicedir} \
	--infifo="" --query-cache --password-cache

# Greylist daemon
${prefix}/sbin/svctool --greylist=1999 --servicedir=${servicedir} --min-resend-min=2 \
    --resend-win-hr=24 --timeout-days=30 --context-file=greylist.context \
    --hash-size=65535 --save-interval=5 --whitelist=greylist.white

if [ ${notlsarr} -eq 0 ] ; then
  # qmail-daned tlsa daemon
  ${prefix}/sbin/qmail-daned > /dev/null 2>&1 || \
  (test $? -ne 100 && ${prefix}/sbin/svctool \
    --tlsa=1998 --servicedir=${servicedir} \
    --timeout-days=30 --context-file=tlsa.context \
    --hash-size=65535 --save-interval=5 --whitelist=tlsa.white ) || \
    echo "Not enabling qmail-daned service" 1>&2
fi

# qmail-qmtpd service
${prefix}/sbin/svctool --qmtp=209 --servicedir=${servicedir} --qbase=${qbase} \
	--qcount=${qcount} --qstart=1 --cntrldir=control --localip=0 \
	--maxdaemons=75 --maxperip=25 \
	--fsync --syncdir --memory=$qmtp_soft_mem --min-free=52428800

# qmail-qmqpd service
${prefix}/sbin/svctool --qmqp=628 --servicedir=${servicedir} --qbase=${qbase} \
	--qcount=${qcount} --qstart=1 --cntrldir=control --localip=0 \
	--maxdaemons=75 --maxperip=25 \
	--fsync --syncdir --memory=$qmqp_soft_mem --min-free=52428800
/bin/touch ${servicedir}/qmail-qmqpd.628/down

# virus/spam filtering
${prefix}/sbin/svctool --qscanq --servicedir=${servicedir} --scanint=200
if [ $have_clamav -eq 1 ] ; then
	${prefix}/sbin/svctool --config=clamd
	${prefix}/sbin/svctool --config=foxhole
	# create clamd, freshclam service
	${prefix}/sbin/svctool --clamd --servicedir=${servicedir} --clamdPrefix=$clamdPrefix \
		--sysconfdir=$mysysconfdir
  if [ $have_clamav -eq 1 ] ; then
    echo "Checking if clamd/freshclam is running"
    count=`ps -e|grep clamd|wc -l`
    if [ $count -gt 0 ] ; then
      echo "Disabling clamd service"
      $TOUCH ${servicedir}/clamd/down
    fi
    count=`ps -e|grep freshclam|wc -l`
    if [ $count -gt 0 ] ; then
      echo "Disabling freshclam service"
      $TOUCH ${servicedir}/freshclam/down
    fi
  fi
fi

# udplogger service
${prefix}/sbin/svctool --udplogger=3000 --localip=0 --timeout=10 --servicedir=${servicedir}

/bin/cat <<EOFUU > ${qsysconfdir}/control/signatures
# Windows executables seen in active virii
TVqQAAMAA
TVpQAAIAA
# Additional windows executable signatures not yet seen in virii
TVpAALQAc
TVpyAXkAX
TVrmAU4AA
TVrhARwAk
TVoFAQUAA
TVoAAAQAA
TVoIARMAA
TVouARsAA
TVrQAT8AA
# .ZIPfile signature seen in SoBig.E and mydoom:
#UEsDBBQAA
#UEsDBAoAAA
# .GIF file found in a previous Microsoft virus making the rounds.
R0lGODlhaAA7APcAAP///+rp6puSp6GZrDUjUUc6Zn53mFJMdbGvvVtXh2xre8bF1x8cU4yLprOy
EOFUU
/bin/chown qscand:qscand ${qsysconfdir}/control/signatures

if [ -x ${sysconfdir}/init.d/apparmor ] ; then
if [ $have_clamav -eq 1 ] ; then
	if [ -d ${sysconfdir}/apparmor.d/local ] ; then
		reload=0
		if [ ! -f ${sysconfdir}/apparmor.d/local/usr.sbin.clamd ] ; then
		(
		echo "# Site-specific additions and overrides for usr.sbin.clamd."
		echo "# For more details, please see ${sysconfdir}/apparmor.d/local/README."
		echo "capability setuid,"
		echo "capability setgid,"
		echo "${shareddir}/clamd/ r,"
		echo "${shareddir}/clamd/** krw,"
		reload=1
		) > ${sysconfdir}/apparmor.d/local/usr.sbin.clamd
		fi
		if [ ! -f ${sysconfdir}/apparmor.d/local/usr.sbin.freshclam ] ; then
		(
		echo "# Site-specific additions and overrides for usr.bin.freshclam."
		echo "# For more details, please see ${sysconfdir}/apparmor.d/local/README."
		echo "${shareddir}/share/clamd/ r,"
		echo "${shareddir}/share/clamd/** krw,"
		reload=1
		) > ${sysconfdir}/apparmor.d/local/usr.bin.freshclam
		fi
		if [ $reload -ne 0 ] ; then
			echo "Reloading apparmor profiles..."
			/usr/sbin/service apparmor reload || true
		fi
	fi
fi
fi

# rebuild cdb
for i in smtp qmtp qmqp
do
	for j in `/bin/ls ${qsysconfdir}/tcp/tcp.$i 2>/dev/null`
	do
		t1=`date +'%s' -r $j`
		if [ -f $j.cdb ] ; then
			t2=`date +'%s' -r $j.cdb`
		else
			t2=0
		fi
		if [ $t1 -gt $t2 ] ; then
			echo "Creating CDB $j.cdb"
			${prefix}/bin/tcprules $j.cdb $j.tmp < $j && /bin/chmod 664 $j.cdb \
				&& /bin/chown indimail:indimail $j.cdb
		fi
	done
done

# turn off automatic refresh for services during first time installation
svc_list=""
for i in clamd greylist.1999 qmail-logfifo qmail-qmqpd.628 \
  qmail-qmtpd.209 qmail-smtpd.25 qmail-smtpd.366 .svscan \
  qmail-smtpd.465 qmail-smtpd.587 qscanq udplogger.3000 \
  qmail-daned.1998
do
  if [ -d ${servicedir}/$i ] ; then
    svc_list="$svc_list ${servicedir}/$i"
  fi
done
svc_list="$svc_list ${qsysconfdir}/control/defaultqueue"
${prefix}/sbin/svctool --servicedir=${servicedir} --norefreshsvc="0 $svc_list"

# libindimail/mysql_lib control file updater
${prefix}/sbin/svctool --libwatch --servicedir=${servicedir}

echo "compressing man pages"
find ${mandir} -type f -exec gzip -q {} \;
if [ -f ${mandir}/man7/authlib.7.gz ] ; then
	for i in authshadow.7 authpwd.7 authpam.7 authcustom.7
	do
		${rm} -f ${mandir}/man7/$i
		ln -sf ${mandir}/man7/authlib.7.gz ${mandir}/man7/$i.gz
	done
fi

#
# Install IndiMail to be started on system boot
# The add-boot command installs svscan to be started by init, systemd, upstart or launchctl
#
echo "adding indimail-mta startup"
if [ /usr/lib/sm.bin/sendmail ] ; then
	if [ -x /usr/sbin/update-alternatives ] ; then
	/usr/sbin/update-alternatives --remove-all sendmail-mta || true
	/usr/sbin/update-alternatives --remove-all sendmail-msp || true
	elif [ -x /usr/bin/update-alternatives ] ; then
	/usr/bin/update-alternatives --remove-all sendmail-mta || true
	/usr/bin/update-alternatives --remove-all sendmail-msp || true
	fi
fi
${prefix}/sbin/svctool --disable-service sendmail exim4
${prefix}/sbin/svctool --config=add-boot
if [ -x /bin/systemctl ] ; then
  /bin/systemctl enable svscan.service
fi

echo "Checking Installed Permissions"
if [ -f ${libexecdir}/instcheck ] ; then
	${libexecdir}/instcheck > /dev/null 2>&1
fi

#
# create /tmp/skip-install_tables so that svctool does not run sbin/install_tables
#
) >> /tmp/indimail-mta-setup.log 2>&1

if [ -f ${sysconfdir}/init/svscan.conf -o -f ${sysconfdir}/event.d/svscan ] ; then
	echo "1. Issue /sbin/initctl emit qmailstart to start services"
	count=1
elif [ -f ${sysconfdir}/systemd/system/multi-user.target.wants/svscan.service ] ; then
	echo "1. Issue /bin/systemctl start svscan to start services"
	count=1
else
	if [ -x ${prefix}/sbin/initsvc ] ; then
	echo "1. Issue ${prefix}/sbin/initsvc -on"
	else
	echo "1. Issue ${sysconfdir}/init.d/svscan start"
	fi
	echo "2. Issue /sbin/init q to start services"
	count=2
fi
default_domain=$([ -n "$HOSTNAME" ] && echo "$HOSTNAME" || `uname -n`)
echo "   You can also use the 'service' command to start indimail-mta"
echo "   sudo service indimail start"
count=`expr $count + 1`
echo "$count. Change your default domain in ${qsysconfdir}/control/defaultdomain"
count=`expr $count + 1`
echo "$count. You can optionally run the following command to verify installation"
echo "   sudo ${prefix}/sbin/svctool --check-install --servicedir=${servicedir} \\"
echo "       --qbase=${qbase} --qcount=${qcount} --qstart=1"
if [ ! -f ${qsysconfdir}/certs/servercert.pem ] ; then
count=`expr $count + 1`
echo "$count. You also need to create CERTS for SSL/TLS. run the following"
echo "   command to create cert"
echo "   ${prefix}/sbin/svctool --config=cert --postmaster=postmaster@${default_domain} --common_name=${default_domain}"
fi
echo
echo "Check /tmp/indimail-mta-setup.log for the detailed installation log!!!"

exit 0
