# daemontools-x/debian/daemontools.postinst.in
#!/bin/sh
# postinst script for indimail
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

prefix=@prefix@
qmaildir=@qmaildir@
mandir=${prefix}/share/man
shareddir=${prefix}/share/indimail
logdir=@logdir@
servicedir=@servicedir@
qsysconfdir=@qsysconfdir@
cp=/bin/cp
rm=/bin/rm
mv=/bin/mv
####################################

setperm()
{
	local user="$1"; shift
	local group="$1"; shift
	local mode="$1"; shift
	local file="$1"; shift
	/usr/sbin/dpkg-statoverride --list "$file" >/dev/null && return 0
	/usr/sbin/dpkg-statoverride --update --add "$user" "$group" "$mode" "$file"
}

fixperms()
{
	if [ $# -eq 3 ] ; then # directory
	#%dir %attr(0555,root,indimail) /var/indimail/sbin
		lin=$2
		file=$3
		if [ ! -d $file ] ; then
			return 0
		fi
	elif [ $# = 2 ] ; then # file
	#%attr(0555,root,qmail) /var/indimail/bin/rsmtp
		lin=$1
		file=$2
		if [ ! -f $file ] ; then
			return 0
		fi
	fi
	perm=`echo $lin|cut -d\( -f2 | cut -d, -f1`
	perm_f1=`echo $perm | cut -c1`
	own=`echo $lin|cut -d\( -f2 | cut -d, -f2`
	grp=`echo $lin|cut -d\( -f2 | cut -d, -f3|cut -d\) -f1`
	if [ -f $file ] ; then
		if [ $perm_f1 -eq 0 ] ; then # 4755, 0644
			/bin/chown $own:$grp $file
			/bin/chmod $perm $file
		else
			echo $file
			if [ -f /usr/sbin/dpkg-statoverride ] ; then
				setperm $own $grp $perm $file
			else
				/bin/chown $own:$grp $file
				/bin/chmod $perm $file
			fi
		fi
	elif [ -d $file ] ; then
		/bin/chown $own:$grp $file
		/bin/chmod $perm $file
	fi
}

ID=`id -u`
if [ $ID -ne 0 ] ; then
	echo "You are not root" >&2
	exit 1
fi

case "$1" in
    configure)
		if [ " $2" = " " ] ; then
		echo "installing daemontools"
		else
		echo "upgrading daemontools $2"
		fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
	exit 0
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

#
# We need to do this for debian. SIGH!!
#
if [ -f ${qsysconfdir}/daemontools_perm_list ] ; then
	echo "Fixing Permissions (debian stupidity)"
(
/bin/grep -v "^#" ${qsysconfdir}/daemontools_perm_list | while read line
do
	if [ " $line" = " " ] ; then
		continue
	fi
	fixperms $line
done
) > ${qsysconfdir}/daemontools-stat.override
fi

if [ -f /tmp/daemontools.upgrade ] ; then # upgrade
  if [ -f ${shareddir}/boot/rpm.init ] ; then
    if [ ! "${prefix}" = "/usr" ] ; then
    echo "recreating ld.so cache"
    /sbin/ldconfig
    fi
    echo "Running Custom Upgrade Script"
    /bin/sh ${shareddir}/boot/rpm.init upgrade
  fi
  if [ -f /etc/init.d/svscan ] ; then
    cmp -s ${prefix}/bin/qmailctl /etc/init.d/svscan >/dev/null 2>&1 || \
      ${cp} ${prefix}/bin/qmailctl /etc/init.d/svscan || true
  fi
  if [ -x /bin/systemctl -a -d /lib/systemd/system ] ; then
    cmp -s ${shareddir}/boot/systemd /lib/systemd/system/svscan.service >/dev/null 2>&1 || true
    ${cp} ${shareddir}/boot/systemd /lib/systemd/system/svscan.service
    /bin/systemctl daemon-reload
  elif [ -x /bin/systemctl -a -d /usr/lib/systemd/system ] ; then
    cmp -s ${shareddir}/boot/systemd /usr/lib/systemd/system/svscan.service >/dev/null 2>&1 || true
    ${cp} ${shareddir}/boot/systemd /usr/lib/systemd/system/svscan.service
    /bin/systemctl daemon-reload
  elif [ -f /sbin/initctl -a -d /etc/init ] ; then
    cmp -s ${shareddir}/boot/upstart /etc/init/svscan.conf >/dev/null 2>&1 || true
    ${cp} ${shareddir}/boot/upstart /etc/init/svscan.conf
  elif [ -d /etc/event.d ] ; then
    cmp -s ${shareddir}/boot/upstart /etc/event.d/upstart >/dev/null 2>&1 || true
    ${cp} ${shareddir}/boot/upstart /etc/event.d
  elif [ -d /System/Library/LaunchDaemons ] ; then
    cmp -s ${shareddir}/boot/indimail.plist /System/Library/LaunchDaemons >/dev/null 2>&1 || true
    ${cp} ${shareddir}/boot/indimail.plist /System/Library/LaunchDaemons
  fi
  ${rm} -f /tmp/daemontools.upgrade
  exit 0
fi
${rm} -f /tmp/daemontools.install

echo "Creating ${logdir}"
if [ ! -d ${logdir} ] ; then
	/bin/mkdir ${logdir}
fi
/bin/chown -R qmaill:nofiles ${logdir}
${prefix}/sbin/minisvc --config=add-boot
if [ -x /bin/systemctl ] ; then
  /bin/systemctl enable svscan.service
fi
# svscanlog service
${prefix}/sbin/minisvc --svscanlog --servicedir=${servicedir} --scanint=60 > /dev/null || true

exit 0
