# indimail-1.9.1/debian/indimail.prerm.in
#!/bin/sh
# prerm script for indimail
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

prefix=@prefix@
indimaildir=@indimaildir@
sysconfdir=@sysconfdir@
shareddir=@shareddir@
mandir=@mandir@
rm=/bin/rm
mv=/bin/mv
nocourierimap=0

ID=`id -u`
if [ $ID -ne 0 ] ; then
	echo "You are not root" >&2
	exit 1
fi

case "$1" in
    remove|upgrade|deconfigure)
	(
	if [ -f ${prefix}/sbin/initsvc ] ; then
		${prefix}/sbin/initsvc -status || /etc/init.d/indimail stop
	fi
	echo "Giving IndiMail exactly 5 seconds to exit nicely" 1>&2
	if test -f /lib/systemd/system/indimail.service && test -x /bin/systemctl
	then
		/bin/systemctl stop indimail.service > /dev/null 2>&1 || true
		/bin/systemctl disable indimail.service > /dev/null 2>&1 || true
		if [ " $1" != " upgrade" ] ; then
			${rm} -f /lib/systemd/system/indimail.service || true
		fi
	elif test -f /etc/init/svscan.conf
	then
		/sbin/initctl emit qmailstop > /dev/null 2>&1 || true
		if [ " $1" != " upgrade" ] ; then
			${rm} -f /etc/event.d/svscan
		fi
	elif test -f /etc/event.d/svscan
	then
		/sbin/initctl emit qmailstop > /dev/null 2>&1 || true
		if [ " $1" != " upgrade" ] ; then
			${rm} -f /etc/event.d/svscan || true
		fi
	elif test -x ${prefix}/sbin/initsvc
	then
		${prefix}/sbin/initsvc -off || true
	fi
	if [ -f ${shareddir}/boot/rpm.init ] ; then
		echo "Running Custom Un-Installation Script for preun" 1>&2
		/bin/sh ${shareddir}/boot/rpm.init preun
	fi
	sleep 5
	) > /tmp/indimail-uninstall.log 2>&1
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

# we are doing upgrade
if [ " $1" = " upgrade" -o " $1" = " deconfigure" ] ; then
	exit 0
fi

if [ -s ${sysconfdir}/stat.override -a -f /usr/sbin/dpkg-statoverride ] ; then
	echo "reverting permissions"
	/bin/cat ${sysconfdir}/stat.override | while read line
	do
		if [ -f $line ] ; then
			echo "reverting permission for $line"
			/usr/sbin/dpkg-statoverride --remove $line || true
		fi
	done
fi
# Remove IndiMail being started on system boot
echo "removing indimail startup"
SERVICE_LIST="sendmail exim4 mysql"
export SERVICE_LIST
$prefix/sbin/svctool --config=rm-boot

if [ $nocourierimap -eq 0 ] ; then
	for i in authshadow.7 authpwd.7 authpam.7 authcustom.7
	do
		${rm} -f ${mandir}/man7/$i.gz || true
	done
	${rmdir} --ignore-fail-on-non-empty ${mandir}/man7 2>/dev/null || true
fi

for i in lib sbin
do
	if [ -f /usr/$i/sendmail.old ] ; then
		${mv} /usr/$i/sendmail.old /usr/$i/sendmail || true
	fi
done
exit 0
