# indimail-mta-x/debian/qmta.postinst.in
#!/bin/sh
# postinst script for qmta
set -e

prefix=@prefix@
mandir=@mandir@
qmaildir=@qmaildir@
qsysconfdir=@qsysconfdir@
sysconfdir=@sysconfdir@
shareddir=@shareddir@
rm=/bin/rm
mv=/bin/mv
cp=/bin/cp

setperm()
{
	local user="$1"; shift
	local group="$1"; shift
	local mode="$1"; shift
	local file="$1"; shift
	/usr/sbin/dpkg-statoverride --list "$file" >/dev/null && return 0
	/usr/sbin/dpkg-statoverride --update --add "$user" "$group" "$mode" "$file"
}

fixperms()
{
	if [ $# -eq 3 ] ; then # directory
		lin=$2
		file=$3
		if [ ! -d $file ] ; then
			return 0
		fi
	elif [ $# = 2 ] ; then # file
		lin=$1
		file=$2
		if [ ! -f $file ] ; then
			return 0
		fi
	fi
	perm=`echo $lin|cut -d\( -f2 | cut -d, -f1`
	perm_f1=`echo $perm | cut -c1`
	own=`echo $lin|cut -d\( -f2 | cut -d, -f2`
	grp=`echo $lin|cut -d\( -f2 | cut -d, -f3|cut -d\) -f1`
	if [ -f $file ] ; then
		if [ $perm_f1 -eq 0 ] ; then # 4755, 0644
			/bin/chown $own:$grp $file
			/bin/chmod $perm $file
		else
			echo $file
			if [ -f /usr/sbin/dpkg-statoverride ] ; then
				setperm $own $grp $perm $file
			else
				/bin/chown $own:$grp $file
				/bin/chmod $perm $file
			fi
		fi
	elif [ -d $file ] ; then
		/bin/chown $own:$grp $file
		/bin/chmod $perm $file
	fi
}

ID=$(id -u)
if [ $ID -ne 0 ] ; then
	echo "You are not root" 1>&2
	exit 1
fi
upgrade=0
case "$1" in
    configure)
		if [ -n "$2" ] ; then
			upgrade=1
		fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
	exit 0
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

if [ -f ${qsysconfdir}/qmta.perms ] ; then
	(
	/bin/grep -v "^#" ${qsysconfdir}/qmta.perms | while read line
	do
		if [ " $line" = " " ] ; then
			continue
		fi
		fixperms $line
	done
	) > ${qsysconfdir}/qmta-stat.override
	if [ ! -s ${qsysconfdir}/qmta-stat.override ] ; then
		/bin/rm -f ${qsysconfdir}/qmta-stat.override
	fi
fi

if [ $upgrade -eq 1 ] ; then # upgrade
  (
  echo "checking if system boot scripts need upgrade"
  if [ -x /bin/systemctl -a -d /lib/systemd/system ] ; then
    cmp -s ${shareddir}/boot/qmta-send.sevice /lib/systemd/system/qmta-send.service >/dev/null 2>&1 || \
      (${cp} ${shareddir}/boot/systemd /lib/systemd/system/qmta-send.service &&
      /bin/systemctl daemon-reload && echo "installed new qmta-send service")
  elif [ -x /bin/systemctl -a -d /usr/lib/systemd/system ] ; then
    cmp -s ${shareddir}/boot/qmta-send.service /usr/lib/systemd/system/qmta-send.service >/dev/null 2>&1 || \
      (${cp} ${shareddir}/boot/systemd /usr/lib/systemd/system/qmta-send.service &&
      /bin/systemctl daemon-reload && echo "installed new qmta-send service")
  fi
  # restart qmta-send service if we find /var/tmp/.qmta-send.down
  if [ -f /var/tmp/.qmta-send.down ] ; then
    ${rm} -f /var/tmp/.qmta-send.down
    echo "Starting qmta-send"
    if test -f ${sysconfdir}/systemd/system/multi-user.target.wants/qmta-send.service
    then
      /bin/systemctl start qmta-send > /dev/null 2>&1
    fi
  fi
  ) >> /var/log/qmta-setup.log 2>&1
  exit 0
fi

(
if [ -f /etc/init.d/sendmail ] ; then
	/etc/init.d/sendmail stop
	if [ -f /sbin/chkconfig ] ; then
		/sbin/chkconfig --del sendmail 2>/dev/null
	fi
fi

if [ -f /etc/pam.d/smtp -a ! -L /etc/pam.d/smtp ] ; then
	echo "! /etc/pam.d/smtp is a file, should be a link" 1>&2
	${mv} /etc/pam.d/smtp /etc/pam.d/smtp.old
fi
for i in mailq newaliases; do
	if [ -f /usr/share/man/man1/$i.1.gz -a ! -L /usr/share/man/man1/$i.1.gz ] ; then
		echo "! /usr/share/man/man1/$i.1.gz is a file, should be a link" 1>&2
		${mv} /usr/share/man/man1/$i.1.gz /usr/share/man/man1/$i.old.1.gz
	fi
done
if [ -x /usr/sbin/update-alternatives ] ; then
	if [ /usr/lib/sm.bin/sendmail ] ; then
		/usr/sbin/update-alternatives --remove-all sendmail-msp || true
	fi
	for i in /usr/lib/sendmail /usr/sbin/sendmail; do
		if [ -f $i -a ! -L $i ]; then
			echo "! $i is a file, should be a link" 1>&2
			${mv} $i $i.old
		fi
	done
	/usr/sbin/update-alternatives --install \
		/usr/sbin/sendmail-msp sendmail-msp ${prefix}/bin/sendmail 120 \
		--slave /usr/share/man/man8/sendmail.8.gz sendmail.8.gz \
			${mandir}/man8/isendmail.8.gz \
		--slave /usr/share/man/man8/sendmail-msp.8.gz sendmail-msp.8.gz \
			${mandir}/man8/isendmail.8.gz \
		--slave /usr/sbin/sendmail sendmail ${prefix}/bin/sendmail \
		--slave /usr/lib/sendmail lib.sendmail ${prefix}/bin/sendmail
else
	for i in /usr/lib/sendmail /usr/sbin/sendmail; do
		if [ -f $i -a ! -L $i ]; then
			echo "! $i is a file, should be a link"
			${mv} $i $i.old
			/bin/ln -s ${prefix}/bin/sendmail $i
		elif [ -L $i ];then
			${mv} $i $i.old
			/bin/ln -s ${prefix}/bin/sendmail $i
		elif [ ! -f $i ];then
			echo "! $i is missing"
			/bin/ln -s ${prefix}/bin/sendmail $i
		fi
	done
fi
if [ ! -f ${qsysconfdir}/control/idhost ] ; then
	default_domain=$([ -n "$HOSTNAME" ] && echo "$HOSTNAME" || uname -n)
	echo $default_domain > ${qsysconfdir}/control/idhost
fi
mkdir -p ${qmaildir}/queue
queue-fix ${qmaildir}/queue/qmta
) >> /var/log/qmta-setup.log 2>&1
exit 0
