# indimail-1.8.1/debian/prerm
#!/bin/sh
# prerm script for indimail
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

prefix=/var/indimail
sysconfdir=/etc
rm=/bin/rm
mv=/bin/mv

case "$1" in
    remove|upgrade|deconfigure)
	if [ -f ${prefix}/sbin/initsvc ] ; then
		${prefix}/sbin/initsvc -status
	fi
	if test -f ${sysconfdir}/event.d/svscan
	then
		/sbin/initctl emit qmailstop > /dev/null 2>&1
		echo "Giving IndiMail exactly 5 seconds to exit nicely"
		sleep 5
		if [ " $1" != " upgrade" ] ; then
			${rm} -f ${sysconfdir}/event.d/svscan
		fi
	elif test -x ${prefix}/sbin/initsvc
	then
		${prefix}/sbin/initsvc -off
	fi
	if [ -f ${prefix}/boot/rpm.init ] ; then
		echo "Running Custom Un-Installation Script for preun"
		/bin/sh ${prefix}/boot/rpm.init preun
	fi
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

# we are doing upgrade
if [ " $1" = " upgrade" ] ; then
	exit 0
fi

# Remove IndiMail being started on system boot
echo "removing indimail startup"
if [ -f ${sysconfdir}/init.d/indimail ] ; then
	if [ -f /sbin/chkconfig ] ; then
		/sbin/chkconfig --del indimail
	fi
	${rm} -f ${sysconfdir}/init.d/indimail
fi
if [ -f ${sysconfdir}/init.d/sendmail ] ; then
	if [ -f /sbin/chkconfig ] ; then
		/sbin/chkconfig --add sendmail
	fi
fi
if [ -x /usr/sbin/alternatives ] ; then
	/usr/sbin/alternatives --remove mta ${prefix}/bin/sendmail
	/usr/sbin/alternatives --auto mta
else
	for i in /usr/lib/sendmail /usr/sbin/sendmail; do
		if [ -f $i.old -o -L $i.old ]; then
			echo "restoring sendmail"
			${rm} -f $i
			${mv} $i.old $i
		fi
	done
fi

exit 0
