# indimail-1.8.9/debian/indimail-mini.postinst.in
#!/bin/sh
# postinst script for indimail-mini
set -e

prefix=@indimaildir@
logdir=@logdir@
sysconfdir=/etc
rm=/bin/rm
mv=/bin/mv

ID=$(id -u)
if [ $ID -ne 0 ] ; then
	echo "You are not root" 1>&2
	exit 1
fi

case "$1" in
    configure)
		if [ " $2" = " " ] ; then
		echo "installing indimail @version@"
		else
		echo "upgrading to indimail $2"
		fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
	exit 0
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

if [ -f /tmp/indimail-devel.upgrade ] ; then # upgrade
	if [ -f ${prefix}/boot/rpm.init ] ; then
		echo "recreating ld.so cache"
		/sbin/ldconfig
		echo "Running Custom Upgrade Script"
		/bin/sh ${prefix}/boot/rpm.init upgrade
	fi
	${rm} -f /tmp/indimail-devel.upgrade
	exit 0
fi
${rm} -f /tmp/indimail-devel.install /tmp/indimail-devel.upgrade
if [ -f ${prefix}/bin/qmail-qmqpc ] ; then
	${rm} -f ${prefix}/bin/qmail-queue
fi
cd ${prefix}/bin
/bin/ln -sf qmail-qmqpc qmail-queue
echo "compressing man pages"
find ${prefix}/man -type f -exec gzip -q {} \;
if [ -f ${sysconfdir}/init.d/sendmail ] ; then
	${sysconfdir}/init.d/sendmail stop
	if [ -f /sbin/chkconfig ] ; then
		/sbin/chkconfig --del sendmail 2>/dev/null
	fi
fi

if [ -x /usr/sbin/update-alternatives ] ; then
	if [ /usr/lib/sm.bin/sendmail ] ; then
		/usr/sbin/update-alternatives --remove-all sendmail-msp || true
	fi
	/usr/sbin/update-alternatives --install \
		/usr/sbin/sendmail-msp sendmail-msp ${prefix}/bin/sendmail 120 \
		--slave /usr/share/man/man8/sendmail.8.gz sendmail.8.gz \
			${prefix}/man/man8/sendmail.8.gz \
		--slave /usr/share/man/man8/sendmail-msp.8.gz sendmail-msp.8.gz \
			${prefix}/man/man8/sendmail.8.gz \
		--slave /usr/sbin/sendmail sendmail ${prefix}/bin/sendmail \
		--slave /usr/lib/sendmail lib.sendmail ${prefix}/bin/sendmail
else
	for i in /usr/lib/sendmail /usr/sbin/sendmail; do
		if [ -f $i -a ! -L $i ]; then
			echo "! $i is a file, should be a link"
			${mv} $i $i.old
			/bin/ln -s ${prefix}/bin/sendmail $i
		elif [ -L $i ];then
			${mv} $i $i.old
			/bin/ln -s ${prefix}/bin/sendmail $i
		elif [ ! -f $i ];then
			echo "! $i is missing"
			/bin/ln -s ${prefix}/bin/sendmail $i
		fi
	done
fi
if [ ! -f ${prefix}/control/idhost ] ; then
	uname -n > ${prefix}/control/idhost
fi
exit 0
