# indimail-mta-2.3/debian/indimail-mini.prerm.in
#!/bin/sh
# prerm script for indimail-mini
#
# see: dh_installdeb(1)

set -e

prefix=@prefix@
shareddir=@shareddir@
rm=/bin/rm
mv=/bin/mv

ID=`id -u`
if [ $ID -ne 0 ] ; then
	echo "You are not root" >&2
	exit 1
fi

case "$1" in
    remove|upgrade|deconfigure)
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

# we are doing upgrade
if [ " $1" = " upgrade" -o " $1" = " deconfigure" ] ; then
	exit 0
fi

# Remove IndiMail being started on system boot
echo "adding sendmail startup"
if [ -f /etc/init.d/sendmail ] ; then
	if [ -f /sbin/chkconfig ] ; then
		/sbin/chkconfig --add sendmail 2>/dev/null
	fi
fi
if [ -x /usr/sbin/alternatives ] ; then
	echo "removing ${prefix}/bin/sendmail as alternative"
	/usr/sbin/alternatives --remove sendmail-msp ${prefix}/bin/sendmail || true
	echo "adding default alternative for sendmail"
	/usr/sbin/alternatives --auto sendmail-msp || true
elif [ -x /usr/sbin/update-alternatives ] ; then
	echo "removing ${prefix}/bin/sendmail as alternative"
	if [ -f /usr/lib/sm.bin/sendmail ] ; then
		/usr/sbin/update-alternatives --remove-all sendmail-msp || true
		echo "adding default alternative for sendmail"
		/usr/sbin/update-alternatives --install \
			/usr/sbin/sendmail-msp sendmail-msp /usr/lib/sm.bin/sendmail 25 \
			--slave /usr/share/man/man8/sendmail-msp.8.gz sendmail-msp.8.gz \
				/usr/share/man/man8/sendmail.sendmail.8.gz \
			--slave /usr/sbin/sendmail sendmail /usr/lib/sm.bin/sendmail \
			--slave /usr/lib/sendmail lib.sendmail /usr/lib/sm.bin/sendmail \
			--slave /usr/share/man/man8/sendmail.8.gz sendmail.8.gz \
				/usr/share/man/man8/sendmail.sendmail.8.gz || true
	else
		/usr/sbin/update-alternatives --remove sendmail-msp ${prefix}/bin/sendmail || true
		echo "adding default alternative for sendmail"
		/usr/sbin/update-alternatives --auto sendmail-msp || true
	fi
fi
for i in /usr/lib/sendmail /usr/sbin/sendmail /etc/pam.d/smtp; do
	if [ -f $i.old -o -L $i.old ]; then
		echo "restoring $i"
		${rm} -f $i
		${mv} $i.old $i
	fi
done
for i in mailq newaliases; do
	if [ -f /usr/share/man/man1/$i.old.1.gz ] ; then
		echo "restoring /usr/share/man/man1/$i.1.gz"
		${mv} /usr/share/man/man1/$i.old.1.gz /usr/share/man/man1/$i.1.gz
	fi
done
exit 0
